# servers/challenger/etc/cron.d/backup-check
# Run backup check at 10am
10 10 * * * root /usr/local/sbin/backup-check
