# etc/rc.d/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local


# set timeout for disks RAID style
#for j in twl0; do
#  for (( i=8; i<24; i++ )); do
#    if smartctl -d 3ware,$i /dev/$j -iq silent; then
#      smartctl -l scterc,70,70 -d 3ware,$i /dev/$j
#    fi
#  done
#done

