# etc/rc.d/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local
if [ -x /etc/cron.daily/nfs-remount ]; then 
  at -f /etc/cron.daily/nfs-remount "now + 15 minutes"
fi
